# cron/magnusrodseth.crontab
0 10 * * 1 /bin/sh ~/dotfiles/scripts/utils/auto-update.sh 2>> ~/dotfiles/scripts/utils/auto-update.log
