# Library/LaunchAgents/com.aoc2025.solver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aoc2025.solver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd /Users/magnusrodseth/dev/personal/aoc-2025 &amp;&amp; /opt/homebrew/bin/claude --dangerously-skip-permissions "/solve"</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>2</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/magnusrodseth/dev/personal/aoc-2025/logs/solver.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/magnusrodseth/dev/personal/aoc-2025/logs/solver.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/magnusrodseth</string>
    </dict>
</dict>
</plist>
