# Library/LaunchAgents/com.aoc2025.solver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aoc2025.solver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/magnusrodseth/dev/personal/aoc-2025/scripts/run-solver.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/magnusrodseth/dev/personal/aoc-2025/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/magnusrodseth/dev/personal/aoc-2025/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/magnusrodseth/dev/personal/aoc-2025</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/magnusrodseth/.cargo/bin</string>
    </dict>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
